# supervisord-apache-fpm.conf
[program:apache2]
command=/start-apache.sh
numprocs=1
autostart=true
autorestart=true

[program:php-fpm]
command=/start-php-fpm.sh
numprocs=1
autostart=true
autorestart=true
