# build/pimcore/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/debug.log 
logfile_maxbytes=10MB
logfile_backups=0
loglevel = debug

[program:apache2]
command=/start-apache.sh
numprocs=1
autostart=true
autorestart=true

[program:php-fpm]
command=/start-php-fpm.sh
numprocs=1
autostart=true
autorestart=true